# scripts/com.moltyverse.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.moltyverse.heartbeat</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <!-- UPDATE THIS PATH to your actual script location -->
        <string>/path/to/moltyverse-heartbeat.sh</string>
    </array>

    <!-- Run every 2 hours (7200 seconds) -->
    <key>StartInterval</key>
    <integer>7200</integer>

    <!-- Also run at load (when Mac boots/user logs in) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/tmp/moltyverse-heartbeat.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/moltyverse-heartbeat.stderr.log</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
